# roles/systemd_startx/templates/startx.service.j2
[Unit]
Description=Start X on login

[Service]
ExecStart=/usr/local/bin/startx-wrapper
Restart=no

[Install]
WantedBy=default.target
