# roles/sshfs_client_managed_user_data/templates/sshfs_mount.service.j2
[Unit]
Description=SSHFS mount to {{ remote_path }}
After=network.target

[Mount]
What={{ remote_path }}
Where={{ local_path }}
Options={{ sshfs_options }}
Type=fuse.sshfs

[Install]
WantedBy=default.target
